# docker/supervisor.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/var/run/supervisord.pid
user=root

; ---------------------------
; PHP-FPM
; ---------------------------
[program:php-fpm]
command=bash -c "docker-php-entrypoint php-fpm"
autostart=true
autorestart=true
startsecs=5
startretries=10
stdout_logfile=/var/log/php-fpm.stdout.log
stderr_logfile=/var/log/php-fpm.stderr.log

; ---------------------------
; Nginx
; ---------------------------
[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/nginx.stdout.log
stderr_logfile=/var/log/nginx.stderr.log

; ---------------------------
; Laravel Queue Worker
; ---------------------------
[program:laravel-worker]
command=php /var/www/html/artisan queue:work --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/laravel-worker.stdout.log
stderr_logfile=/var/log/laravel-worker.stderr.log

; ---------------------------
; Laravel Scheduler
; ---------------------------
[program:laravel-scheduler]
command=php /var/www/html/artisan schedule:work
autostart=true
autorestart=true
stdout_logfile=/var/log/laravel-scheduler.stdout.log
stderr_logfile=/var/log/laravel-scheduler.stderr.log

; ---------------------------
; Laravel Echo Server
; ---------------------------
; [program:laravel-echo-server]
; command=laravel-echo-server start --config=/var/www/html/laravel-echo-server.json
; directory=/var/www/html
; autostart=true
; autorestart=true
stdout_logfile=/var/log/laravel-echo-server.stdout.log
stderr_logfile=/var/log/laravel-echo-server.stderr.log